# src/backups/motion_detectiond.original
#!/bin/bash
### BEGIN INIT INFO
# Provides: MotionDetection.py
# Required-Start:    $local_fs $syslog $remote_fs dbus
# Required-Stop:     $local_fs $syslog $remote_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start MotionDetection system on boot
### END INIT INFO

/usr/bin/python /home/pi/Documents/Python/MotionDetection/MotionDetection.py -e 'sshmonitorapp@gmail.com' -p 'GHOST21ghost' &
